# deploy/tablet/dropbear-loopback.socket
[Unit]
Description=Loopback-only Dropbear socket for private Paperboard forwarding
Conflicts=dropbear.service

[Socket]
ListenStream=127.0.0.1:2222
Accept=yes

[Install]
WantedBy=sockets.target
Also=dropbearkey.service
