# deploy/tablet/paperboard-tailscale-serve.service
[Unit]
Description=Private Paperboard SSH forwarding over Tailscale Serve
After=paperboard-tailscale.service dropbear-loopback.socket
Requires=paperboard-tailscale.service dropbear-loopback.socket

[Service]
Type=oneshot
ExecStart=/home/root/.local/share/paperboard/tailscale/current/tailscale --socket=/home/root/.local/share/paperboard/tailscale/runtime/tailscaled.sock serve --bg --yes --tcp=22 tcp://127.0.0.1:2222
RemainAfterExit=yes
Restart=on-failure
RestartSec=15
TimeoutStartSec=30
NoNewPrivileges=yes
PrivateTmp=yes

[Install]
WantedBy=multi-user.target
